# packaging/systemd/netfoil.socket
[Unit]
Description=netfoil socket

[Socket]
ListenDatagram=127.0.0.1:53
Service=netfoil.service

[Install]
WantedBy=sockets.target
